# etc/supervisor.d/gunicorn.ini
[program:DBW301]
command=gunicorn -c etc/gunicorn.conf.py main:app --reload
autostart=true
autorestart=true
numprocs=1
priority = 900
redirect_stderr=true
stderr_logfile=/var/log/%(program_name)s.log
stdout_logfile=/var/log/%(program_name)s.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=100
stdout_logfile_backups=10
